# workspace/crontab/laradock
* * * * * laradock php /var/www/rtc/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock php /var/www/app2018/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock php /var/www/iosclub/artisan schedule:run >> /dev/null 2>&1

